# rpi/roles/monitoring/files/hc-ping.timer
[Unit]
Description=Ping healthchecks.io hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
