# rpi/roles/monitoring/templates/hc-ping.service.j2
[Service]
Type=oneshot
ExecStart=/usr/bin/curl -fsS -m 10 --retry 5 -o /dev/null \
  https://hc-ping.com/{{ hcping_uuid['content'] | b64decode }}
# vim: ft=systemd.jinja2
